# contrib/vagrant/setup/dotfiles/dot_setup-desktop.conf
description "Setup Desktop"
start on desktop-start
task
script
  if [ -e ~/.setup-desktop.sh ]; then
    bash ~/.setup-desktop.sh
    rm -f ~/.setup-desktop.sh
  fi
end script
